# services/predictor/predictor.service
[Unit]
Description=CogniFlight Predictor Service
After=network.target
Wants=network.target

[Service]
Type=notify
User=jeremia
Group=jeremia
WorkingDirectory=/home/jeremia/projects/cogniflight-edge-backup/services/predictor
Environment=PATH=/home/jeremia/projects/cogniflight-edge-backup/services/predictor/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/jeremia/projects/cogniflight-edge-backup
ExecStart=/home/jeremia/projects/cogniflight-edge-backup/services/predictor/.venv/bin/python main.py
Restart=on-failure
WatchdogSec=30
StartLimitBurst=3
StartLimitIntervalSec=60
RestartSec=5
StandardOutput=journal
StandardError=journal


# Resource limits - Optimized for RPi5 16GB (Critical fatigue prediction)
MemoryMax=2G
CPUQuota=75%
CPUWeight=900
Nice=-10

# Security settings
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target