# services/alert_manager/alert_manager.service
[Unit]
Description=CogniFlight Alert Manager Service
After=network.target
Wants=network.target

[Service]
Type=notify
User=jeremia
Group=jeremia
SupplementaryGroups=gpio i2c spi
WorkingDirectory=/home/jeremia/projects/cogniflight-edge-backup/services/alert_manager
Environment=PATH=/home/jeremia/projects/cogniflight-edge-backup/services/alert_manager/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/jeremia/projects/cogniflight-edge-backup
ExecStart=/home/jeremia/projects/cogniflight-edge-backup/services/alert_manager/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
WatchdogSec=30
StartLimitBurst=3
StartLimitIntervalSec=60
StandardOutput=journal
StandardError=journal

# Resource limits - Optimized for RPi5 16GB (Critical safety alerts)
MemoryMax=1.5G
CPUQuota=50%
CPUWeight=700
Nice=-5

# Security settings
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target