# services/face_recognition/face_recognition.service
# /home/jeremia/Desktop/cogniflight-edge/services/face_recognition/face_recognition.service

[Unit]
Description=CogniFlight Face Recognition Service
After=network.target alert-manager.service
Wants=network.target
Requires=alert_manager.service

[Service]
Type=simple
User=jeremia
Group=jeremia
WorkingDirectory=/home/jeremia/Desktop/cogniflight-edge/services/face_recognition
Environment=PATH=/home/jeremia/Desktop/cogniflight-edge/services/face_recognition/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HEARTBEAT_DIR=/run/cogniflight
ExecStart=/home/jeremia/Desktop/cogniflight-edge/services/face_recognition/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Runtime directory for heartbeat files
RuntimeDirectory=cogniflight
RuntimeDirectoryMode=0755

# Resource limits - Optimized for RPi5 16GB
MemoryMax=3G
CPUQuota=50%
CPUWeight=600
Nice=5

# Security settings
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target