# services/face_recognition/face_recognition.service
[Unit]
Description=CogniFlight Face Recognition Service
After=network.target
Wants=network.target

[Service]
Type=notify
User=jeremia
Group=jeremia
WorkingDirectory=/home/jeremia/projects/cogniflight-edge-backup/services/face_recognition
Environment=PATH=/home/jeremia/projects/cogniflight-edge-backup/services/face_recognition/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/jeremia/projects/cogniflight-edge-backup
ExecStart=/home/jeremia/projects/cogniflight-edge-backup/services/face_recognition/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
WatchdogSec=30
StartLimitBurst=3
StartLimitIntervalSec=60
StandardOutput=journal
StandardError=journal

# Resource limits - Optimized for RPi5 16GB (AI inference workload)
MemoryMax=4G
CPUQuota=75%
CPUWeight=800
Nice=0

# Security settings
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target